# do-lister-bot.service
[Unit]
Description=Do Lister Telegram Bot
After=network.target

[Service]
Type=simple
User=user1
WorkingDirectory=/home/user1/do-lister
ExecStart=/usr/bin/python3 /home/user1/do-lister/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
